# systemd/idigest-digest.service
[Unit]
Description=idigest weekly digest email
Requires=idigest-llm.service
After=idigest-llm.service

[Service]
Type=oneshot
WorkingDirectory=/home/sur/repo/local-inference
ExecStart=/home/sur/repo/local-inference/.venv/bin/idigest digest
