# systemd/idigest-digest.timer
[Unit]
Description=idigest weekly digest (Sunday 08:00)

[Timer]
OnCalendar=Sun *-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
