# systemd/idigest-email.timer
[Unit]
Description=idigest daily email at 07:30 (content pre-generated at 05:00)

[Timer]
OnCalendar=*-*-* 07:30:00
Persistent=true

[Install]
WantedBy=timers.target
